# conf/templates/crontab
0 23 * * * cd {{ BASE_DIR }} && make regular_backup && make media_backup